# deploy/systemd/gov-backend.service
[Unit]
Description=Government Project Backend Service
After=network.target mariadb.service

[Service]
Type=simple
User=gov-backend
Group=gov-backend
WorkingDirectory=/opt/peoplesAirDefence

# Environment variables file
EnvironmentFile=/opt/peoplesAirDefence/etc/gov-backend.env

# Default profile (override via environment variable)
Environment="SPRING_PROFILES_ACTIVE=prod"

# Start command
ExecStart=/opt/peoplesAirDefence/bin/start.sh

# Restart policy
Restart=on-failure
RestartSec=10

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gov-backend

# Resource limits
LimitNOFILE=65535
LimitNPROC=65535

# Timeout settings
TimeoutStartSec=300
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
